# couchdb-livesync/rootfs/etc/s6-overlay/s6-rc.d/couchdb/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: CouchDB for Obsidian LiveSync
# Runs CouchDB with configuration for Obsidian LiveSync
# ==============================================================================

# Source bashio library
source /usr/lib/bashio/bashio.sh

declare username
declare password
declare database

# Read configuration from options.json
if [[ -f /data/options.json ]]; then
    username=$(bashio::config 'username')
    password=$(bashio::config 'password')
    database=$(bashio::config 'database')
else
    echo "Warning: /data/options.json not found, using defaults"
    username="admin"
    password="admin"
    database="obsidian-livesync"
fi

echo "Starting CouchDB for Obsidian LiveSync..."

# Ensure data directory exists
mkdir -p /data/couchdb
chown -R couchdb:couchdb /data/couchdb

# Link data directory to CouchDB default location
if [[ ! -L /opt/couchdb/data ]]; then
    rm -rf /opt/couchdb/data
    ln -s /data/couchdb /opt/couchdb/data
fi

# Configure CouchDB local.ini with admin credentials and Obsidian LiveSync settings
cat > /opt/couchdb/etc/local.d/obsidian-livesync.ini << EOF
[couchdb]
single_node = true
max_document_size = 50000000

[admins]
${username} = ${password}

[chttpd]
port = 5984
bind_address = 0.0.0.0
require_valid_user = true
enable_cors = true
max_http_request_size = 4294967296

[chttpd_auth]
require_valid_user = true

[httpd]
WWW-Authenticate = Basic realm="couchdb"
enable_cors = true

[cors]
credentials = true
origins = app://obsidian.md,capacitor://localhost,http://localhost

[log]
level = info
EOF

# Set correct ownership
chown couchdb:couchdb /opt/couchdb/etc/local.d/obsidian-livesync.ini

echo "CouchDB configured with user: ${username}"
echo "Database: ${database}"
echo "Listening on port 5984"

# Start CouchDB (using the official entrypoint)
exec /opt/couchdb/bin/couchdb
