# couchdb-livesync/rootfs/etc/services.d/couchdb/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: CouchDB for Obsidian LiveSync
# Takes down the S6 supervisor when CouchDB fails
# ==============================================================================

declare exit_code=${1}

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.error "CouchDB failed with exit code ${exit_code}"
    bashio::log.warning "Halting add-on..."
    /run/s6/basedir/bin/halt
fi
