# couchdb-livesync/rootfs/etc/services.d/couchdb/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: CouchDB for Obsidian LiveSync
# Runs CouchDB with configuration for Obsidian LiveSync
# ==============================================================================

declare username
declare password
declare database

# Read configuration
username=$(bashio::config 'username')
password=$(bashio::config 'password')
database=$(bashio::config 'database')

bashio::log.info "Starting CouchDB for Obsidian LiveSync..."

# Ensure data directory exists and has correct permissions
mkdir -p /data/couchdb
chown -R couchdb:couchdb /data/couchdb

# Configure CouchDB local.ini with admin credentials and Obsidian LiveSync settings
cat > /etc/couchdb/local.ini << EOF
[couchdb]
single_node = true
max_document_size = 50000000

[admins]
${username} = ${password}

[chttpd]
port = 5984
bind_address = 0.0.0.0
require_valid_user = true
enable_cors = true
max_http_request_size = 4294967296

[chttpd_auth]
require_valid_user = true

[httpd]
WWW-Authenticate = Basic realm="couchdb"
enable_cors = true

[cors]
credentials = true
origins = app://obsidian.md,capacitor://localhost,http://localhost

[log]
level = info
EOF

# Set correct ownership
chown couchdb:couchdb /etc/couchdb/local.ini

bashio::log.info "CouchDB configured with user: ${username}"
bashio::log.info "Database: ${database}"
bashio::log.info "Listening on port 5984"

# Start CouchDB
exec /usr/bin/couchdb
